# shorewall/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-r-xr--r--	root/root	etc/rc.d/shorewall
drwxr-xr-x	root/root	etc/shorewall/
-rw-------	root/root	etc/shorewall/Makefile
-rw-------	root/root	etc/shorewall/accounting
-rw-r--r--	root/root	etc/shorewall/actions
-rw-------	root/root	etc/shorewall/blacklist
-rw-------	root/root	etc/shorewall/continue
-rw-------	root/root	etc/shorewall/ecn
-rw-------	root/root	etc/shorewall/hosts
-rw-------	root/root	etc/shorewall/init
-rw-------	root/root	etc/shorewall/initdone
-rw-------	root/root	etc/shorewall/interfaces
-rw-------	root/root	etc/shorewall/ipsec
-rw-------	root/root	etc/shorewall/maclist
-rw-------	root/root	etc/shorewall/masq
-rw-------	root/root	etc/shorewall/nat
-rw-------	root/root	etc/shorewall/netmap
-rw-r--r--	root/root	etc/shorewall/params
-rw-------	root/root	etc/shorewall/policy
-rw-------	root/root	etc/shorewall/providers
-rw-------	root/root	etc/shorewall/proxyarp
-rw-------	root/root	etc/shorewall/route_rules
-rw-------	root/root	etc/shorewall/routestopped
-rw-------	root/root	etc/shorewall/rules
-rw-r--r--	root/root	etc/shorewall/shorewall.conf
-rw-------	root/root	etc/shorewall/start
-rw-------	root/root	etc/shorewall/started
-rw-------	root/root	etc/shorewall/stop
-rw-------	root/root	etc/shorewall/stopped
-rw-------	root/root	etc/shorewall/tcclasses
-rw-------	root/root	etc/shorewall/tcdevices
-rw-------	root/root	etc/shorewall/tcrules
-rw-------	root/root	etc/shorewall/tos
-rw-------	root/root	etc/shorewall/tunnels
-rwxr--r--	root/root	etc/shorewall/zones
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/shorewall
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/shorewall-accounting.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-actions.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-blacklist.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-ecn.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-exclusion.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-hosts.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-interfaces.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-maclist.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-masq.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-modules.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-nat.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-nesting.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-netmap.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-params.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-policy.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-providers.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-proxyarp.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-rfc1918.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-route_rules.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-routestopped.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-rules.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-tcclasses.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-tcdevices.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-tcrules.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-tos.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-tunnels.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-vardir.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall-zones.5.gz
-rw-r--r--	root/root	usr/man/man5/shorewall.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/shorewall.8.gz
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/shorewall-perl/
drwxr-xr-x	root/root	usr/share/shorewall-perl/Shorewall/
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Accounting.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Actions.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Chains.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Compiler.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Config.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/IPAddrs.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Nat.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Policy.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Proc.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Providers.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Proxyarp.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Rules.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Tc.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Tunnels.pm
-rw-r--r--	root/root	usr/share/shorewall-perl/Shorewall/Zones.pm
-rwxr-xr-x	root/root	usr/share/shorewall-perl/compiler.pl
-rw-r--r--	root/root	usr/share/shorewall-perl/prog.footer
-rw-r--r--	root/root	usr/share/shorewall-perl/prog.functions
-rw-r--r--	root/root	usr/share/shorewall-perl/prog.header
-rw-r--r--	root/root	usr/share/shorewall-perl/version
drwxr-xr-x	root/root	usr/share/shorewall/
-rw-r--r--	root/root	usr/share/shorewall/action.Drop
-rw-r--r--	root/root	usr/share/shorewall/action.Reject
-rw-r--r--	root/root	usr/share/shorewall/action.template
-rw-r--r--	root/root	usr/share/shorewall/actions.std
drwxr-xr-x	root/root	usr/share/shorewall/configfiles/
-rw-r--r--	root/root	usr/share/shorewall/configfiles/Makefile
-rw-r--r--	root/root	usr/share/shorewall/configfiles/accounting
-rw-r--r--	root/root	usr/share/shorewall/configfiles/actions
-rw-r--r--	root/root	usr/share/shorewall/configfiles/blacklist
-rw-r--r--	root/root	usr/share/shorewall/configfiles/continue
-rw-r--r--	root/root	usr/share/shorewall/configfiles/ecn
-rw-r--r--	root/root	usr/share/shorewall/configfiles/hosts
-rw-r--r--	root/root	usr/share/shorewall/configfiles/init
-rw-r--r--	root/root	usr/share/shorewall/configfiles/initdone
-rw-r--r--	root/root	usr/share/shorewall/configfiles/interfaces
-rw-r--r--	root/root	usr/share/shorewall/configfiles/ipsec
-rw-r--r--	root/root	usr/share/shorewall/configfiles/maclist
-rw-r--r--	root/root	usr/share/shorewall/configfiles/masq
-rw-r--r--	root/root	usr/share/shorewall/configfiles/nat
-rw-r--r--	root/root	usr/share/shorewall/configfiles/netmap
-rw-r--r--	root/root	usr/share/shorewall/configfiles/params
-rw-r--r--	root/root	usr/share/shorewall/configfiles/policy
-rw-r--r--	root/root	usr/share/shorewall/configfiles/providers
-rw-r--r--	root/root	usr/share/shorewall/configfiles/proxyarp
-rw-r--r--	root/root	usr/share/shorewall/configfiles/route_rules
-rw-r--r--	root/root	usr/share/shorewall/configfiles/routestopped
-rw-r--r--	root/root	usr/share/shorewall/configfiles/rules
-rw-r--r--	root/root	usr/share/shorewall/configfiles/shorewall.conf
-rw-r--r--	root/root	usr/share/shorewall/configfiles/start
-rw-r--r--	root/root	usr/share/shorewall/configfiles/started
-rw-r--r--	root/root	usr/share/shorewall/configfiles/stop
-rw-r--r--	root/root	usr/share/shorewall/configfiles/stopped
-rw-r--r--	root/root	usr/share/shorewall/configfiles/tcclasses
-rw-r--r--	root/root	usr/share/shorewall/configfiles/tcdevices
-rw-r--r--	root/root	usr/share/shorewall/configfiles/tcrules
-rw-r--r--	root/root	usr/share/shorewall/configfiles/tos
-rw-r--r--	root/root	usr/share/shorewall/configfiles/tunnels
-rw-r--r--	root/root	usr/share/shorewall/configfiles/zones
-rw-r--r--	root/root	usr/share/shorewall/configpath
-rwxr-xr-x	root/root	usr/share/shorewall/firewall
lrwxrwxrwx	root/root	usr/share/shorewall/functions -> lib.base
-rw-r--r--	root/root	usr/share/shorewall/lib.base
-rw-r--r--	root/root	usr/share/shorewall/lib.cli
-rw-r--r--	root/root	usr/share/shorewall/lib.config
-rw-r--r--	root/root	usr/share/shorewall/lib.dynamiczones
-rw-r--r--	root/root	usr/share/shorewall/macro.AllowICMPs
-rw-r--r--	root/root	usr/share/shorewall/macro.Amanda
-rw-r--r--	root/root	usr/share/shorewall/macro.Auth
-rw-r--r--	root/root	usr/share/shorewall/macro.BitTorrent
-rw-r--r--	root/root	usr/share/shorewall/macro.CVS
-rw-r--r--	root/root	usr/share/shorewall/macro.DNS
-rw-r--r--	root/root	usr/share/shorewall/macro.Distcc
-rw-r--r--	root/root	usr/share/shorewall/macro.Drop
-rw-r--r--	root/root	usr/share/shorewall/macro.DropDNSrep
-rw-r--r--	root/root	usr/share/shorewall/macro.DropUPnP
-rw-r--r--	root/root	usr/share/shorewall/macro.Edonkey
-rw-r--r--	root/root	usr/share/shorewall/macro.FTP
-rw-r--r--	root/root	usr/share/shorewall/macro.Finger
-rw-r--r--	root/root	usr/share/shorewall/macro.GRE
-rw-r--r--	root/root	usr/share/shorewall/macro.Gnutella
-rw-r--r--	root/root	usr/share/shorewall/macro.HTTP
-rw-r--r--	root/root	usr/share/shorewall/macro.HTTPS
-rw-r--r--	root/root	usr/share/shorewall/macro.ICQ
-rw-r--r--	root/root	usr/share/shorewall/macro.IMAP
-rw-r--r--	root/root	usr/share/shorewall/macro.IMAPS
-rw-r--r--	root/root	usr/share/shorewall/macro.IPIP
-rw-r--r--	root/root	usr/share/shorewall/macro.IPP
-rw-r--r--	root/root	usr/share/shorewall/macro.IPPserver
-rw-r--r--	root/root	usr/share/shorewall/macro.IPsec
-rw-r--r--	root/root	usr/share/shorewall/macro.IPsecah
-rw-r--r--	root/root	usr/share/shorewall/macro.IPsecnat
-rw-r--r--	root/root	usr/share/shorewall/macro.JabberPlain
-rw-r--r--	root/root	usr/share/shorewall/macro.JabberSecure
-rw-r--r--	root/root	usr/share/shorewall/macro.Jabberd
-rw-r--r--	root/root	usr/share/shorewall/macro.Jetdirect
-rw-r--r--	root/root	usr/share/shorewall/macro.L2TP
-rw-r--r--	root/root	usr/share/shorewall/macro.LDAP
-rw-r--r--	root/root	usr/share/shorewall/macro.LDAPS
-rw-r--r--	root/root	usr/share/shorewall/macro.MySQL
-rw-r--r--	root/root	usr/share/shorewall/macro.NNTP
-rw-r--r--	root/root	usr/share/shorewall/macro.NNTPS
-rw-r--r--	root/root	usr/share/shorewall/macro.NTP
-rw-r--r--	root/root	usr/share/shorewall/macro.NTPbrd
-rw-r--r--	root/root	usr/share/shorewall/macro.PCA
-rw-r--r--	root/root	usr/share/shorewall/macro.POP3
-rw-r--r--	root/root	usr/share/shorewall/macro.POP3S
-rw-r--r--	root/root	usr/share/shorewall/macro.Ping
-rw-r--r--	root/root	usr/share/shorewall/macro.PostgreSQL
-rw-r--r--	root/root	usr/share/shorewall/macro.Printer
-rw-r--r--	root/root	usr/share/shorewall/macro.RDP
-rw-r--r--	root/root	usr/share/shorewall/macro.RNDC
-rw-r--r--	root/root	usr/share/shorewall/macro.Rdate
-rw-r--r--	root/root	usr/share/shorewall/macro.Reject
-rw-r--r--	root/root	usr/share/shorewall/macro.Rsync
-rw-r--r--	root/root	usr/share/shorewall/macro.SMB
-rw-r--r--	root/root	usr/share/shorewall/macro.SMBBI
-rw-r--r--	root/root	usr/share/shorewall/macro.SMBswat
-rw-r--r--	root/root	usr/share/shorewall/macro.SMTP
-rw-r--r--	root/root	usr/share/shorewall/macro.SMTPS
-rw-r--r--	root/root	usr/share/shorewall/macro.SNMP
-rw-r--r--	root/root	usr/share/shorewall/macro.SPAMD
-rw-r--r--	root/root	usr/share/shorewall/macro.SSH
-rw-r--r--	root/root	usr/share/shorewall/macro.SVN
-rw-r--r--	root/root	usr/share/shorewall/macro.SixXS
-rw-r--r--	root/root	usr/share/shorewall/macro.Submission
-rw-r--r--	root/root	usr/share/shorewall/macro.Syslog
-rw-r--r--	root/root	usr/share/shorewall/macro.TFTP
-rw-r--r--	root/root	usr/share/shorewall/macro.Telnet
-rw-r--r--	root/root	usr/share/shorewall/macro.Telnets
-rw-r--r--	root/root	usr/share/shorewall/macro.Time
-rw-r--r--	root/root	usr/share/shorewall/macro.Trcrt
-rw-r--r--	root/root	usr/share/shorewall/macro.VNC
-rw-r--r--	root/root	usr/share/shorewall/macro.VNCL
-rw-r--r--	root/root	usr/share/shorewall/macro.Web
-rw-r--r--	root/root	usr/share/shorewall/macro.Webmin
-rw-r--r--	root/root	usr/share/shorewall/macro.Whois
-rw-r--r--	root/root	usr/share/shorewall/macro.template
-rw-------	root/root	usr/share/shorewall/modules
-rw-r--r--	root/root	usr/share/shorewall/rfc1918
-rw-r--r--	root/root	usr/share/shorewall/version
-rwxr-xr-x	root/root	usr/share/shorewall/wait4ifup
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/shorewall/
